# workspace/systemd/openclaw-tool-ain-phi.service
[Unit]
Description=OpenClaw Tool - ain.phi (127.0.0.1:18990)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
Type=simple
WorkingDirectory=%h/src/clawd
Environment=OPENCLAW_NODE_ID=dali
ExecStart=%h/src/clawd/workspace/scripts/run_tool_ain_phi.sh
Restart=always
RestartSec=3
EnvironmentFile=-%h/.config/openclaw/tools.env
StandardOutput=append:%h/.local/state/openclaw/tools/ain.phi.out.log
StandardError=append:%h/.local/state/openclaw/tools/ain.phi.err.log

[Install]
WantedBy=default.target
