# workspace/systemd/openclaw-user-inference-distill.timer
[Unit]
Description=Periodic OpenClaw user inference distillation

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
Unit=openclaw-user-inference-distill.service

[Install]
WantedBy=timers.target
