# workspace/systemd/openclaw-idle-reaper.service
[Unit]
Description=OpenClaw Idle Reaper (stop GPU tools when idle)

[Service]
Type=oneshot
WorkingDirectory=%h/src/clawd
ExecStart=%h/src/clawd/workspace/scripts/idle_reaper.py
