# workspace/systemd/openclaw-contract-manager.timer
[Unit]
Description=Run OpenClaw Contract Manager every 2 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=120
AccuracySec=10

[Install]
WantedBy=timers.target
