# workspace/systemd/openclaw-idle-reaper.timer
[Unit]
Description=Run OpenClaw Idle Reaper every 2 minutes

[Timer]
OnBootSec=45
OnUnitActiveSec=120
AccuracySec=10

[Install]
WantedBy=timers.target
