# workspace/systemd/openclaw-discord-bridge.service
[Unit]
Description=OpenClaw Discord webhook bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/src/clawd
ExecStart=%h/.local/share/openclaw/venvs/discord-bot/bin/python %h/src/clawd/workspace/scripts/discord_webhook_bridge.py --env-file %h/.config/openclaw/discord-bridge.env
StandardOutput=append:%h/.local/state/openclaw/discord/bridge.out.log
StandardError=append:%h/.local/state/openclaw/discord/bridge.err.log

