# workspace/systemd/openclaw-discord-bridge.service
[Unit]
Description=OpenClaw Discord webhook bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/src/clawd
ExecStart=/usr/bin/python3 %h/src/clawd/workspace/scripts/discord_webhook_bridge.py --env-file %h/.config/openclaw/discord-bridge.env
StandardOutput=journal
StandardError=journal
