# workspace/systemd/openclaw-discord-bridge.timer
[Unit]
Description=Run OpenClaw Discord webhook bridge every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=openclaw-discord-bridge.service

[Install]
WantedBy=timers.target
