# workspace/systemd/openclaw-itc-telegram.service
[Unit]
Description=OpenClaw ITC Telegram listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/src/clawd
EnvironmentFile=-%h/.config/openclaw/itc-telegram.env
ExecStartPre=/usr/bin/mkdir -p %h/.local/state/openclaw
ExecStart=/usr/bin/env bash %h/src/clawd/scripts/run_itc_telegram_service.sh
Restart=on-failure
RestartSec=10
StandardOutput=append:%h/.local/state/openclaw/itc-telegram.log
StandardError=append:%h/.local/state/openclaw/itc-telegram.log

[Install]
WantedBy=default.target
