# workspace/systemd/openclaw-contract-manager.service
[Unit]
Description=OpenClaw Contract Manager (dynamic SERVICE/CODE)
After=default.target

[Service]
Type=oneshot
WorkingDirectory=%h/src/clawd
ExecStart=%h/src/clawd/workspace/scripts/contract_manager.py
