# workspace/launchd/ai.openclaw.market-sentiment-idle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.market-sentiment-idle</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/heathyeager/clawd/workspace/scripts/market_sentiment_idle_runner.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/heathyeager/Library/Logs/openclaw-market-sentiment-idle.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/heathyeager/Library/Logs/openclaw-market-sentiment-idle.log</string>
</dict>
</plist>
