# workspace/launchd/ai.openclaw.weekly-x-strategy-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.weekly-x-strategy-review</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/heathyeager/clawd/workspace/scripts/run_weekly_x_strategy_review.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>18</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/heathyeager/Library/Logs/openclaw-weekly-x-strategy-review.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/heathyeager/Library/Logs/openclaw-weekly-x-strategy-review.log</string>
</dict>
</plist>
