# workspace/source-ui/source-ui.service
[Unit]
Description=Source UI - OpenClaw System Management Dashboard (user unit template)
After=network.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/openclaw/source-ui.env
Environment=SOURCE_UI_HOST=127.0.0.1
Environment=SOURCE_UI_PORT=18990
Environment=GATEWAY_URL=http://127.0.0.1:18789
WorkingDirectory=%h/src/clawd/workspace/source-ui
ExecStart=/usr/bin/env python3 %h/src/clawd/workspace/source-ui/app.py --host ${SOURCE_UI_HOST} --port ${SOURCE_UI_PORT}
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=source-ui

[Install]
WantedBy=default.target
